# systemd/pmg-quarantine-action-server.service
[Unit]
Description=PMG Quarantine Junk — action server
Documentation=https://github.com/mapazzzm/pmg-quarantine-junk
After=network.target postgresql.service
Wants=network.target

[Service]
Type=simple
User=pmg-quarantine
Group=pmg-quarantine
ExecStart=/usr/local/bin/pmg-quarantine-action-server
Restart=on-failure
RestartSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=pmg-quarantine-action-server

# Ограничения безопасности
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/lib/pmg-quarantine-junk /var/log /run
ReadOnlyPaths=/etc/pmg-quarantine-junk /etc/pmg /usr/local/lib/pmg-quarantine-junk

[Install]
WantedBy=multi-user.target
